# config/runit/freshclam/run
#!/bin/sh
exec freshclam -d --foreground=true | logger -t freshclam
