# config/runit/amavis/finish
#!/bin/sh
rm -f /run/amavis/amavisd.pid
