# config/runit/clamd/run
#!/bin/sh
exec /usr/sbin/clamd --foreground=true | logger -t clamd
